# debian/postinst
#!/bin/sh

update-rc.d wb-mqtt-astra defaults
