# debian/postinst
#!/bin/sh

update-rc.d wb-mqtt-astra defaults
invoke-rc.d wb-mqtt-astra status >/dev/null 2>&1 && invoke-rc.d wb-mqtt-astra restart; true
